# conductor-agent.service
[Unit]
Description=Conductor agent

[Service]
Type=simple
ExecStart=/root/conductor-agent/start.sh
Restart=no
